# etc/periodic/daily/invoiceninja
#!/bin/sh
/var/www/artisan ninja:send-invoices
/var/www/artisan ninja:send-reminders
